# roles/metrics/templates/graphite.service
[Service]
ExecStart=/usr/bin/graphite-web-manage runserver

[Install]
WantedBy=default.target
