# roles/metrics/templates/statsd.service
[Service]
ExecStart=/usr/bin/statsdpy-service -f --conf=/etc/statsd.conf

[Install]
WantedBy=default.target
